# java_test_deb_deb_maven_plugin/src/deb/java-test-deb.service
[Unit]
Description=test daemon

[Service]
ExecStart=/usr/bin/java -jar /opt/java-test-deb/java-test-deb-1.deb

[Install]
WantedBy=multi-user.target
